# cmd/supervisor.aiadmin.conf
[program:aiadmin]
command=/home/pinfire/pkg/aiadmin/aiadmin
directory=/tmp/
autostart=true
autorestart=true
redirect_stderr=true
user=root